# roles/gws_sync/templates/aurora-mirror-verify.service.j2
[Unit]
Description=Aurora mirror verification and prune-gate manifests
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{ aurora_service_user }}
Group={{ aurora_service_group }}
ExecStart=/usr/local/bin/aurora-mirror-verify
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
