# roles/dashboard_services/templates/aurora-mobile-api.service.j2
[Unit]
Description=Aurora Dashboard mobile API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ aurora_service_user }}
Group={{ aurora_service_group }}
WorkingDirectory={{ aurora_app_dir }}
EnvironmentFile=/etc/aurora-dashboard.env
ExecStart={{ aurora_venv }}/bin/uvicorn mobile_api:app --host {{ aurora_mobile_api_host }} --port {{ aurora_mobile_api_port }}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
