# roles/dashboard_services/templates/aurora-power-soc-forecast.service.j2
[Unit]
Description=Aurora power SOC forecast from ECMWF solar data
After=network-online.target aurora-power-append.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ aurora_service_user }}
Group={{ aurora_service_group }}
WorkingDirectory={{ aurora_app_dir }}
EnvironmentFile=/etc/aurora-dashboard.env
ExecStart={{ aurora_venv }}/bin/python generate_power_soc_forecast.py --power-zarr {{ aurora_zarr.power }} --output-zarr {{ aurora_zarr.power_soc_forecast }} --cache-dir {{ aurora_paths.products_power_ecmwf_solar_forecast }} --state {{ aurora_power_forecast_root }}/power_soc_forecast_state.json --archive-zarr {{ aurora_zarr.power_soc_forecast_archive }} --skill-zarr {{ aurora_zarr.power_soc_forecast_skill }} --hindcast-zarr {{ aurora_zarr.power_soc_hindcast }} --horizon-hours 96
Nice=10
TimeoutStartSec=900
