# roles/dashboard_services/templates/aurora-power-operating-scenarios.timer.j2
[Unit]
Description=Fallback update for development learned APS operating plans

[Timer]
OnBootSec=3min
OnUnitInactiveSec=10min
RandomizedDelaySec=30
Unit=aurora-power-operating-scenarios.service

[Install]
WantedBy=timers.target
