# roles/dashboard_services/templates/aurora-dashboard-health-probe.service.j2
[Unit]
Description=Measure Aurora development and production dashboard health
After=network-online.target aurora-dashboard.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ aurora_service_user }}
Group={{ aurora_service_group }}
ExecStart=/usr/local/bin/aurora-dashboard-health-probe
Nice=15
CPUQuota=10%
TimeoutStartSec=90
