# roles/dashboard_services/templates/aurora-pdu-append.service.j2
[Unit]
Description=Aurora PDU Zarr append
After=network-online.target aurora-pdu-source-sync.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ aurora_service_user }}
Group={{ aurora_service_group }}
WorkingDirectory={{ aurora_app_dir }}
EnvironmentFile=/etc/aurora-dashboard.env
ExecStart={{ aurora_venv }}/bin/python append_new_pdu_to_zarr.py --root {{ aurora_paths.raw_pdu }} --zarr {{ aurora_zarr.pdu }} --chunk-time 1440
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
