# roles/dashboard_services/templates/aurora-dashboard.service.j2
[Unit]
Description=Aurora Cloud Panel Dashboard
After=network-online.target
Wants=network-online.target

[Service]
User={{ aurora_service_user }}
Group={{ aurora_service_group }}
WorkingDirectory={{ aurora_app_dir }}
EnvironmentFile=/etc/aurora-dashboard.env
ExecStart={{ aurora_venv }}/bin/panel serve app.py model-evaluation.py dart.py \
  --address={{ aurora_panel_host }} \
  --port={{ aurora_panel_port }} \
{% for origin in aurora_dashboard_allowed_origins %}
  --allow-websocket-origin={{ origin }} \
{% endfor %}
  --static-dirs dashboard-assets={{ aurora_app_dir }}/assets wxcam-media=/data/aurora/products/wxcam auroracam-media={{ aurora_paths.raw_auroracam }} \
  --keep-alive=15000 \
  --check-unused-sessions={{ aurora_dashboard_session_check_ms }} \
  --unused-session-lifetime={{ aurora_dashboard_unused_session_lifetime_ms }} \
  --session-token-expiration=86400 \
  --num-procs={{ aurora_dashboard_num_procs }}
Restart=always
RestartSec=5
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
