# roles/dashboard_services/templates/aurora-ecmwf-provider-gate.service.j2
[Unit]
Description=Evaluate the Earthkit decoder shadow promotion gate
After=aurora-ecmwf-provider-shadow.service

[Service]
Type=oneshot
User={{ aurora_service_user }}
Group={{ aurora_service_group }}
WorkingDirectory={{ aurora_app_dir }}
EnvironmentFile=/etc/aurora-dashboard.env
ExecStart={{ aurora_venv }}/bin/python evaluate_ecmwf_provider_gate.py --history {{ aurora_ecmwf_shadow_history_path }} --output {{ aurora_ecmwf_shadow_gate_path }} --minimum-days 7 --minimum-samples 50 --maximum-ssrd-difference 0.001
Nice=15
CPUQuota=20%
