# roles/dashboard_services/templates/aurora-dev-power-display-refresh.service.j2
[Unit]
Description=Refresh compact Power display products from mirrored Power data
After=network-online.target aurora-dev-live-pull@product-power.service
Wants=network-online.target
ConditionPathExists={{ aurora_zarr.power }}

[Service]
Type=oneshot
User={{ aurora_service_user }}
Group={{ aurora_service_group }}
WorkingDirectory={{ aurora_app_dir }}
EnvironmentFile=/etc/aurora-dashboard.env
ExecStart=/usr/bin/timeout --signal=TERM --kill-after=30s 10m /bin/bash -lc '{{ aurora_venv }}/bin/python generate_power_display_summary.py && {{ aurora_venv }}/bin/python generate_power_quicklooks.py --prewarm-only --skip-display-refresh'
Nice=15
CPUQuota=100%
CPUWeight=10
MemoryHigh=3G
MemoryMax=4G
TimeoutStartSec=11min
