# roles/dashboard_services/templates/aurora-power-soc-forecast-learn.timer.j2
[Unit]
Description=Fallback fast SOC and operating-mode forecast re-anchor every 10 minutes

[Timer]
OnUnitInactiveSec=10min
RandomizedDelaySec=30s
Persistent=true
Unit=aurora-power-soc-forecast-learn.service

[Install]
WantedBy=timers.target
