# roles/dashboard_services/templates/aurora-power-operating-scenarios.timer.j2
[Unit]
Description=Update development learned APS operating plans every 30 minutes

[Timer]
OnBootSec=3min
OnUnitInactiveSec=30min
RandomizedDelaySec=30
Unit=aurora-power-operating-scenarios.service

[Install]
WantedBy=timers.target
