# roles/standby_replication/templates/aurora-standby-pull.service.j2
[Unit]
Description=Aurora warm-standby pull replication
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/bin/aurora-standby-pull
TimeoutStartSec={{ aurora_standby_replication_timeout_start_sec }}
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6

