# roles/source_sync/templates/aurora-auroracam-backfill.service.j2
[Unit]
Description=Aurora newest-first AURORACam archive backfill over Tailscale SSH
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ aurora_service_user }}
Group={{ aurora_service_group }}
ExecStart=/usr/local/bin/aurora-auroracam-backfill
TimeoutStartSec=10min
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
