# roles/source_sync/templates/aurora-radar-backfill.service.j2
[Unit]
Description=Aurora bounded parallel native cloud-radar archive backfill
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ aurora_service_user }}
Group={{ aurora_service_group }}
ExecStart=/usr/local/bin/aurora-radar-backfill
TimeoutStartSec=12h
Nice=0
IOSchedulingClass=best-effort
IOSchedulingPriority=2
