# roles/source_sync/templates/aurora-uas-source-sync.service.j2
[Unit]
Description=Aurora pull UAS Menapia MQTT log over Tailscale SSH
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ aurora_service_user }}
Group={{ aurora_service_group }}
ExecStart=/usr/local/bin/aurora-uas-sync
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
